# root/etc/services.d/flood/run
#!/usr/bin/with-contenv sh

cd /flood
yarn start